# supervisord.conf
; =============================================================================
; GusEngine V10 — Supervisord Configuration
; Manages: Qdrant, BGE-M3 Embedding Server, FastAPI Backend
;
; Usage: supervisord -c /workspace/GusEngine/supervisord.conf
; Control: supervisorctl -c /workspace/GusEngine/supervisord.conf status
; =============================================================================

[supervisord]
nodaemon=true
logfile=/workspace/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/workspace/supervisord.pid

; --- Unix socket for supervisorctl ---
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; Service 1: Qdrant Vector Database
; Priority 10 = starts first
; =============================================================================
[program:qdrant]
command=/workspace/qdrant --storage-path /workspace/GusEngine/storage/qdrant
directory=/tmp
priority=10
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/workspace/qdrant.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

; =============================================================================
; Service 2: BGE-M3 Embedding Server
; Priority 20 = starts after Qdrant
; startsecs=30 = model takes ~30s to load; don't mark "started" until then
; =============================================================================
[program:embedding]
command=/workspace/venvs/embed/bin/python3 bge_m3_server.py
directory=/workspace/GusEngine
environment=BGE_M3_MODEL_PATH="/workspace/GusEngine/storage/models/bge-m3",TEI_PORT="8080"
priority=20
autostart=true
autorestart=true
startretries=3
startsecs=30
stopwaitsecs=15
stdout_logfile=/workspace/tei.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

; =============================================================================
; Service 3: FastAPI Backend (Gemini 2.5 Flash)
; Priority 30 = starts last (after Qdrant + Embedding)
; ALL environment variables are declared here — never need to set them manually
; =============================================================================
[program:backend]
command=/workspace/venvs/backend/bin/python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8888
directory=/workspace/GusEngine
environment=
    VLLM_BASE_URL="https://generativelanguage.googleapis.com/v1beta/openai",
    VLLM_MODEL="gemini-2.5-flash",
    LLM_API_KEY="%(ENV_LLM_API_KEY)s",
    QDRANT_URL="http://localhost:6333",
    TEI_BASE_URL="http://localhost:8080",
    EMBEDDING_URL="http://localhost:8080",
    TOKENIZER_MODEL_PATH="/workspace/GusEngine/storage/models/Qwen2.5-32B-Instruct-AWQ",
    SYSTEM_PROMPT_PATH="/workspace/GusEngine/config/system_prompt.txt",
    SYSTEM_PROMPT_TOKENS="1800",
    MAX_CONTEXT_TOKENS="32768",
    RESPONSE_BUDGET_TOKENS="8192",
    LEDGER_MAX_TOKENS="2550",
    ALLOWED_PDF_DIR="/workspace/GusEngine/storage/pdfs",
    LEDGER_PATH="/workspace/GusEngine/config/MASTER_LEDGER.md",
    FAILURE_MANIFEST_PATH="/workspace/GusEngine/storage/extracted/.ingest_failures.log",
    EASYOCR_MODULE_PATH="/workspace/GusEngine/storage/easyocr_models",
    DOCLING_MODELS_PATH="/workspace/GusEngine/storage/models/docling-models",
    HF_HUB_OFFLINE="1",
    TRANSFORMERS_OFFLINE="1",
    PYTHONPATH="/workspace/GusEngine"
priority=30
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/workspace/backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
